# debian/hwclock.rules
# This file causes the hardware clock to be set when /dev/rtc is available.

KERNEL=="rtc", ACTION=="add", RUN+="set_hwclock"
